# jukebox.service
[Unit]
Description=Raspberry Pi Jukebox Application
Documentation=https://github.com/kaninfod/jukebox
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=pi
Group=gpio
SupplementaryGroups=spi i2c gpio
WorkingDirectory=/home/pi/shared/jukebox
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/home/pi/shared/jukebox:/home/pi/shared/jukebox/app
Environment=VIRTUAL_ENV=/home/pi/shared/jukebox/venv
Environment=PATH=/home/pi/shared/jukebox/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/pi/shared/jukebox/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --log-level info --proxy-headers --forwarded-allow-ips="*"
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jukebox

# Security settings (modified to allow system operations like reboot/shutdown)
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/home/pi/shared/jukebox
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
